# shadowd.service
[Unit]
Description=Secure login distribution service

[Service]
ExecStart=/usr/bin/shadowd
Restart=always

[Install]
WantedBy=multi-user.target
